# contrib/init/pexad.conf
description "Pexa Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env PEXAD_BIN="/usr/bin/pexad"
env PEXAD_USER="pexa"
env PEXAD_GROUP="pexa"
env PEXAD_PIDDIR="/var/run/pexad"
# upstart can't handle variables constructed with other variables
env PEXAD_PIDFILE="/var/run/pexad/pexad.pid"
env PEXAD_CONFIGFILE="/etc/pexa/pexa.conf"
env PEXAD_DATADIR="/var/lib/pexad"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # pexad will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$PEXAD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run pexad."
        echo "The setting must appear in $PEXAD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Pexa Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$PEXAD_PIDDIR"
    chmod 0755 "$PEXAD_PIDDIR"
    chown $PEXAD_USER:$PEXAD_GROUP "$PEXAD_PIDDIR"
    chown $PEXAD_USER:$PEXAD_GROUP "$PEXAD_CONFIGFILE"
    chmod 0660 "$PEXAD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$PEXAD_PIDFILE" \
    --chuid $PEXAD_USER:$PEXAD_GROUP \
    --exec "$PEXAD_BIN" \
    -- \
    -pid="$PEXAD_PIDFILE" \
    -conf="$PEXAD_CONFIGFILE" \
    -datadir="$PEXAD_DATADIR" \
    -disablewallet \
    -daemon

